# ukrainian-bot.service
# 🧠😂🔥 Systemd сервіс для україномовного Telegram-бота 🧠😂🔥
#
# Встановлення:
# 1. Скопіюйте цей файл до /etc/systemd/system/ukrainian-bot.service
# 2. sudo systemctl daemon-reload
# 3. sudo systemctl enable ukrainian-bot
# 4. sudo systemctl start ukrainian-bot
#
# Управління:
# sudo systemctl status ukrainian-bot
# sudo systemctl restart ukrainian-bot
# sudo systemctl stop ukrainian-bot
# sudo journalctl -u ukrainian-bot -f

[Unit]
Description=🧠😂🔥 Україномовний Telegram-бот з мемами та анекдотами
Documentation=https://github.com/yourusername/ukrainian-telegram-bot
After=network-online.target
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=botuser
Group=botuser
WorkingDirectory=/opt/ukrainian-bot
Environment=PATH=/opt/ukrainian-bot/venv/bin
ExecStart=/opt/ukrainian-bot/venv/bin/python main.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=ukrainian-bot

# Безпека
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ukrainian-bot/logs /opt/ukrainian-bot/data /opt/ukrainian-bot/backups
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Ресурси
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=1G
CPUQuota=100%

# Змінні середовища (альтернативно можна використовувати EnvironmentFile)
Environment=PYTHONPATH=/opt/ukrainian-bot
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONUNBUFFERED=1
Environment=LOG_LEVEL=INFO

# Файл з змінними середовища
EnvironmentFile=-/opt/ukrainian-bot/.env

# Watchdog (опціонально)
WatchdogSec=30

# Автоматичний restart при збоях
RestartPreventExitStatus=SIGKILL
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

# ===================================================================
# ДОДАТКОВІ КОНФІГУРАЦІЇ
# ===================================================================

# Для development версії
# [Unit]
# Description=🧠😂🔥 Україномовний Telegram-бот (Development)
# After=network-online.target
# 
# [Service]
# Type=simple
# User=developer
# WorkingDirectory=/home/developer/ukrainian-bot
# ExecStart=/home/developer/ukrainian-bot/venv/bin/python main.py --dev
# Environment=DEBUG=True
# Environment=LOG_LEVEL=DEBUG
# Restart=on-failure
# RestartSec=5
# 
# [Install]
# WantedBy=multi-user.target

# ===================================================================
# КОМАНДИ ДЛЯ НАЛАШТУВАННЯ
# ===================================================================

# 1. Створення користувача системи:
# sudo adduser --system --group --home /opt/ukrainian-bot botuser

# 2. Встановлення бота:
# sudo mkdir -p /opt/ukrainian-bot
# sudo chown botuser:botuser /opt/ukrainian-bot
# sudo -u botuser git clone https://github.com/yourusername/ukrainian-bot.git /opt/ukrainian-bot
# cd /opt/ukrainian-bot
# sudo -u botuser python -m venv venv
# sudo -u botuser venv/bin/pip install -r requirements.txt

# 3. Конфігурація:
# sudo -u botuser cp .env.example .env
# sudo -u botuser nano .env

# 4. Створення директорій:
# sudo -u botuser mkdir -p /opt/ukrainian-bot/{logs,data,backups}

# 5. Встановлення сервісу:
# sudo cp scripts/ukrainian-bot.service /etc/systemd/system/
# sudo systemctl daemon-reload
# sudo systemctl enable ukrainian-bot

# 6. Запуск:
# sudo systemctl start ukrainian-bot

# 7. Перевірка статусу:
# sudo systemctl status ukrainian-bot
# sudo journalctl -u ukrainian-bot -f

# 8. Автоматичні backup (додати до crontab):
# 0 2 * * * /opt/ukrainian-bot/venv/bin/python /opt/ukrainian-bot/scripts/backup.py create

# 9. Ротація логів (/etc/logrotate.d/ukrainian-bot):
# /opt/ukrainian-bot/logs/*.log {
#     daily
#     rotate 30
#     compress
#     delaycompress
#     missingok
#     notifempty
#     copytruncate
#     su botuser botuser
# }